# CMake/install/clio.service.in
[Unit]
Description=Clio XRPL API server
Documentation=https://github.com/XRPLF/clio.git

After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@CLIO_INSTALL_DIR@/bin/clio_server @CLIO_INSTALL_DIR@/etc/config.json
Restart=on-failure
User=clio
Group=clio
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
