# CMake/clio.service.in
[Unit]
Description=Clio XRPL API server
Documentation=https://github.com/cjcobb23/clio
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@CLIO_INSTALL_DIR@/bin/clio_server @CLIO_INSTALL_DIR@/etc/config.json
Restart=on-failure
User=clio
Group=clio

[Install]
WantedBy=multi-user.target
